# timmerjob/timerjob.cron
0 0 * * * /bin/bash /alidata/www/default/timmerjob/dbbackup.sh 
