# timmerjob/timerjob.cron
0 0 * * * /bin/bash /alidata/timmerjob/dbbackup.sh 
